# p21Home_v11/p21Home_v11/etc/init/proto21_svc.service
[Unit]
Description=uWSGI proto21_home server instance
After=syslog.target

[Service]
ExecStart=/usr/local/bin/uwsgi --ini-paste /webapps/proto21/p21Home/webapps/proto21_home/production.ini
# Requires systemd version 211 or newer
RuntimeDirectory=/webapps/proto21/p21Home/webapps/proto21_home/
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target